# .scripts/vartmp
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /etc/init.d/vartmp, v2.0 2012/02/14 -tclover Exp $

conf_file="/etc/conf.d/${SVCNAME}"
description="use tmpfs and tarball archive to restore/save /var/{log,...}"
:	${tmp_dir:=/run}

depend() {
        need localmount 
        before bootmisc acpid
}

start() {
	if [[ -z "$(mount -t tmpfs | grep ${tmp_dir})" ]]; then 
		mkdir -p "${tmp_dir}" 
		mount tmpdir ${tmp_dir} -t tmpfs -o rw,nosuid,nodev,relatime,mode=755
		eend "$?"; fi
	for dir in ${var_dir}; do check_tar ${dir}; done
    eend "$?"
}

stop() { for dir in ${var_dir}; do pack_dir ${dir}; done; }

check_tar() {
	if [[ -f /var/.${1}.tgz ]]; then unpack_dir ${1}
	elif [[ -d "/var/${1}" ]]; then 
		ewarn "packing /var/${1} for the 1st time"
		cd /var && tar czpf .${1}.tgz ${1}/
		eend "$?" "failed to pack /var/.${1}.tgz"
		mv ${1} "${tmp_dir}"/ && ln -s "${tmp_dir}"/${1} ${1}
		eend "$?" "failed to symlink /var/${1} to  ${tmp_dir}/${1}"
	else mkdir -p "${tmp_dir}"/${1} && chmod 775 "${tmp_dir}"/${1}; fi
}

unpack_dir() {
	ebegin "unpacking /var/.${1}.tgz to ${tmp_dir}/${1}"
	tar xzpf /var/.${1}.tgz -C "${tmp_dir}"/
	eend "$?" "failed to unpack /var/.${1}.tgz"
}

pack_dir() {
	ebegin "packing ${tmp_dir}/${1}"
	cd "${tmp_dir}" && tar czpf /var/.${1}.tgz ${1}/
	eend "$?" "failed to pack /var/.${1}.tgz"
}

# vim:fenc=utf-8:ft=gentoo-init-d:ci:pi:sts=0:sw=4:ts=4:
